# marask7/monthly-maintenance.service
[Unit]
Description=Monthly system update/upgrade/cleanup for Tachyon validator

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update_upgrade_clean.sh
Nice=10
